# 0x1A-application_server/gunicorn.service
[Unit]
Description=Gunicorn instance to serve web application
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/path/to/your/application
ExecStart=/path/to/virtualenv/bin/gunicorn \
          --workers 3 \
          --bind 0.0.0.0:5003 \
          --access-logfile /tmp/airbnb-access.log \
          --error-logfile /tmp/airbnb-error.log \
          your_module_name:app

[Install]
WantedBy=multi-user.target
